# roles/monitoring/templates/check-sm-state.timer.j2
[Unit]
Description=Run check_sm_state.sh every 30 sec

[Timer]
OnBootSec=10
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
