# roles/monitoring/templates/check-sm-state.service.j2
[Unit]
Description=Check Infiniband Subnet Manager State

[Service]
Type=oneshot
ExecStart=/usr/local/bin/check_sm_state.sh
