# roles/monitoring/templates/cleanup-metrics.service.j2
[Unit]
Description=Cleanup stale Prometheus node exporter metrics

[Service]
Type=oneshot
ExecStart=/usr/local/bin/cleanup-metrics.sh
