# roles/monitoring/templates/gpu-process-exporter.service.j2
[Unit]
Description=GPU Process Metrics Exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/gpu-process-exporter.sh
