# roles/monitoring/templates/iblinkinfo-exporter.service.j2
[Unit]
Description=Iblinkinfo Exporter Metrics Exporter

[Service]
Type=oneshot
ExecStart=/usr/local/bin/ib_sw_exp.py
