# roles/monitoring/templates/iblinkinfo-exporter.timer.j2
[Unit]
Description=Run Iblinkinfo Exporter on SWITCH every 30 sec

[Timer]
OnBootSec=10
OnUnitActiveSec=30

[Install]
WantedBy=timers.target
